# config/etc/supervisor/conf.d/apache2.conf
[program:apache2]
priority=100
command=/usr/sbin/apache2ctl -DFOREGROUND