# config/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/sbin/php-fpm8.1 -F
autostart=true
autorestart=unexpected
priority=100
stdout_events_enabled=true
stderr_events_enabled=true